# templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic

[Timer]
OnCalendar={{ borg__borgmatic_schedule }}

#RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
