# templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic

[Timer]
OnCalendar={{ borg__borgmatic_timer_oncalendar }}
RandomizedDelaySec={{ borg__borgmatic_timer_randomizeddelayseconds }}

[Install]
WantedBy=timers.target
